# scheduler/systemd/radio-silence.service
[Unit]
Description=Radio Silence Stream (fallback for /stream)
After=network.target icecast2.service
Wants=icecast2.service

[Service]
Type=simple
ExecStart=/usr/bin/ffmpeg -re -f lavfi -i anullsrc=r=44100:cl=stereo -acodec libmp3lame -ab 192k -f mp3 -content_type audio/mpeg icecast://source:${SOURCE_PASSWORD}@localhost:8000/silence
Restart=always
RestartSec=5

# Environment - override in /etc/radio-scheduler.env
EnvironmentFile=-/etc/radio-scheduler.env

[Install]
WantedBy=multi-user.target
